# ctf1.crontab
* * * * * /usr/bin/flock -n /tmp/runner1.lockfile /usr/bin/python /home/runner1/ctf1.py
